# scripts/sentinelflux.service
[Unit]
Description=SentinelFlux Dashboard
After=network.target

[Service]
User=sentinelflux
WorkingDirectory=/home/sentinelflux/app
EnvironmentFile=/home/sentinelflux/app/.env
ExecStart=/home/sentinelflux/app/.venv/bin/python -m uvicorn dashboard.app:app \
    --host 127.0.0.1 --port 8765 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
